# systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnUnitInactiveSec=5m
Unit=offlineimap.service

[Install]
WantedBy=mail.target
