# systemd/user/offlineimap.service
[Unit]
Description=Offlineimap Service

[Service]
Type=oneshot
Environment=
ExecStart=/usr/bin/offlineimap -o

[Install]
WantedBy=mail.target
