# systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnUnitInactiveSec=1min
OnBootSec=1min
Unit=offlineimap.service

[Install]
WantedBy=mail.target
